# .docker/supervisord.ini
[supervisord]
nodaemon=true
user=root
[program:nginx]
command=nginx
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:php-fpm]
command=/etc/init.d/php8.2-fpm start
stdout_logfile=/dev/stdout
startretries=10
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0